# VCAC-A/Intel_Media_Analytics_Node/scripts/build_scripts/persistent_files/91-netup.rules
SUBSYSTEM!="net", GOTO="net_end"
ACTION!="add", GOTO="net_end"
ENV{INTERFACE}=="eth0", RUN+="/sbin/vca_up.sh", GOTO="net_end"
NAME=="eth0", RUN+="/sbin/vca_up.sh", GOTO="net_end"
LABEL="net_end"
