# apps/rules/97-host_eth_up.rules
SUBSYSTEM=="net", ACTION=="add", ATTR{address}=="fe:00:00:00:0?:0?", KERNEL=="eth*", RUN+="/usr/sbin/vca_eth_ifup.sh '%E{DEVPATH}'"
