# apps/rules/99-vca_mgr_extd.rules
ACTION=="add", KERNEL=="vca_mgr_extd?", GROUP="vcausers", MODE="0660", RUN+="/etc/vca_config.d/vca_auto_boot.sh scan %n"